# supervisord.conf
[supervisord]
nodaemon=true

[program:mlflow]
command=python scripts/train_model.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/mlflow.log
stderr_logfile=/var/log/mlflow_err.log

[program:streamlit]
command=streamlit run app/app.py --server.port 8501 --server.enableCORS false --server.headless true
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/streamlit.log
stderr_logfile=/var/log/streamlit_err.log
